# server/fable-votes.service
[Unit]
Description=Vote API for whatisfabledoing.com
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/fable-votes/votes_server.py
WorkingDirectory=/opt/fable-votes
User=fable-votes
Group=fable-votes
Restart=on-failure
RestartSec=2

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/fable-votes
PrivateTmp=true

[Install]
WantedBy=multi-user.target
